# .config/systemd/user/wm.target
[Unit]
Description=Window manager target
Requires=dbus.socket
AllowIsolate=true

[Install]
Alias=default.target
